# infra/remote/cost-monitor/cost-monitor.timer
[Unit]
Description=snoringcat-platform daily cost monitor (09:00 UTC)

[Timer]
OnCalendar=*-*-* 09:00:00 UTC
Persistent=true
Unit=cost-monitor.service

[Install]
WantedBy=timers.target
